# ansible/files/etc/cron.daily/letsencrypt
#!/bin/sh

systemctl stop apache2
certbot renew --preferred-challenges http > /var/log/letsencrypt/renew.log
systemctl start apache2
